# packages/systemd/ndpisrvd.service
[Unit]
Description=nDPIsrvd Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/nDPIsrvd -c /var/run/ndpisrvd-collector
ExecStopPost=/bin/rm -f /var/run/ndpisrvd-collector
Restart=on-failure

[Install]
WantedBy=multi-user.target
